# lib/systemd/system/avahi-daemon.service.d/10-configure-pi-top-hostname.conf
[Service]
ExecStartPost=/usr/bin/pt-hostname publish
