# lib/systemd/system/systemd-hostnamed.service.d/10-configure-hostname.conf
[Service]
ExecStartPre=/usr/bin/set-pt-hostname
