# src/70-wlan-disable-powersave.rules
SUBSYSTEM=="net", ENV{DEVTYPE}=="wlan", ACTION=="add", RUN+="/sbin/iw dev $name set power_save off"
